# deploy/systemd/mcpx.service.d/hardening.conf
# mcpx hardening drop-in.
#
# DELIBERATELY PARTIAL. mcpx exists to run arbitrary commands and reach other
# hosts (see docs/adr/0002-no-sandbox-full-access.md). The directives that
# would actually confine it - ProtectSystem=strict, NoNewPrivileges,
# ProtectHome, PrivateDevices, a restrictive SystemCallFilter - would break
# the shell/SSH/PTY/sudo capability this service is for, while implying a
# containment guarantee it cannot honour. Those are intentionally OMITTED.
#
# What is applied here is the safe envelope: resource caps so a runaway tool
# cannot exhaust the host, a private /tmp, and the Protect* directives that do
# not interfere with command execution. The real boundary is the service
# account and its credentials, plus the audit/policy layer and the network
# edge - not this drop-in.

[Service]
# Resource envelope (tune to the host).
MemoryHigh=768M
MemoryMax=1024M
CPUQuota=80%
TasksMax=128
LimitNOFILE=8192

# Safe isolation that does not break execution.
PrivateTmp=true
ProtectClock=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
RestrictRealtime=true
LockPersonality=true
UMask=0077

# Encrypted credential delivery (systemd-creds). Create with:
#   systemd-creds encrypt plaintext /etc/credstore.encrypted/mcpx-oauth ...
# then reference the decrypted path via $CREDENTIALS_DIRECTORY in mcpx.env.
# LoadCredentialEncrypted=mcpx-oauth:/etc/credstore.encrypted/mcpx-oauth
