# deploy/systemd/relay-shell.service
[Unit]
Description=relay-shell - MCP server for shell and SSH operations
Documentation=https://github.com/rmednitzer/relay-shell
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=relay-shell
Group=relay-shell
WorkingDirectory=/var/lib/relay-shell
# Configuration is environment-driven; keep secrets out of the unit.
EnvironmentFile=-/etc/relay-shell/relay-shell.env
Environment=RELAY_SHELL_TRANSPORT=http
Environment=RELAY_SHELL_HTTP_HOST=127.0.0.1
Environment=RELAY_SHELL_HTTP_PORT=8080
Environment=RELAY_SHELL_AUDIT_PATH=/var/log/relay-shell/audit.jsonl
ExecStart=/var/lib/relay-shell/venv/bin/relay-shell
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=300
StandardOutput=journal
StandardError=journal
SyslogIdentifier=relay-shell

[Install]
WantedBy=multi-user.target
